# pi-hole-server-5.9-3/files
%FILES%
etc/
etc/dnsmasq.conf
etc/dnsmasq.d/
etc/dnsmasq.d/01-pihole.conf
etc/pihole/
etc/pihole/adlists.list
etc/pihole/dns-servers.conf
etc/pihole/logrotate
etc/sudoers.d/
etc/sudoers.d/pihole
opt/
opt/pihole/
opt/pihole/COL_TABLE
opt/pihole/basic-install.sh
opt/pihole/chronometer.sh
opt/pihole/database_migration/
opt/pihole/database_migration/gravity-db.sh
opt/pihole/database_migration/gravity/
opt/pihole/database_migration/gravity/10_to_11.sql
opt/pihole/database_migration/gravity/11_to_12.sql
opt/pihole/database_migration/gravity/12_to_13.sql
opt/pihole/database_migration/gravity/13_to_14.sql
opt/pihole/database_migration/gravity/14_to_15.sql
opt/pihole/database_migration/gravity/1_to_2.sql
opt/pihole/database_migration/gravity/2_to_3.sql
opt/pihole/database_migration/gravity/3_to_4.sql
opt/pihole/database_migration/gravity/4_to_5.sql
opt/pihole/database_migration/gravity/5_to_6.sql
opt/pihole/database_migration/gravity/6_to_7.sql
opt/pihole/database_migration/gravity/7_to_8.sql
opt/pihole/database_migration/gravity/8_to_9.sql
opt/pihole/database_migration/gravity/9_to_10.sql
opt/pihole/gravity.db.sql
opt/pihole/gravity.sh
opt/pihole/gravity_copy.sql
opt/pihole/list.sh
opt/pihole/mimic_setupVars.conf.sh
opt/pihole/pihole-reenable.sh
opt/pihole/piholeARPTable.sh
opt/pihole/piholeDebug.sh
opt/pihole/piholeLogFlush.sh
opt/pihole/query.sh
opt/pihole/utils.sh
opt/pihole/version.sh
opt/pihole/webpage.sh
srv/
srv/http/
srv/http/pihole/
srv/http/pihole/admin/
srv/http/pihole/admin/api.php
srv/http/pihole/admin/api_FTL.php
srv/http/pihole/admin/api_db.php
srv/http/pihole/admin/auditlog.php
srv/http/pihole/admin/cname_records.php
srv/http/pihole/admin/composer.json
srv/http/pihole/admin/composer.lock
srv/http/pihole/admin/db_graph.php
srv/http/pihole/admin/db_lists.php
srv/http/pihole/admin/db_queries.php
srv/http/pihole/admin/debug.php
srv/http/pihole/admin/dns_records.php
srv/http/pihole/admin/gravity.php
srv/http/pihole/admin/groups-adlists.php
srv/http/pihole/admin/groups-clients.php
srv/http/pihole/admin/groups-domains.php
srv/http/pihole/admin/groups.php
srv/http/pihole/admin/img/
srv/http/pihole/admin/img/boxed-bg-dark.jpg
srv/http/pihole/admin/img/boxed-bg.jpg
srv/http/pihole/admin/img/donate.gif
srv/http/pihole/admin/img/favicons/
srv/http/pihole/admin/img/favicons/android-chrome-192x192.png
srv/http/pihole/admin/img/favicons/android-chrome-512x512.png
srv/http/pihole/admin/img/favicons/apple-touch-icon.png
srv/http/pihole/admin/img/favicons/favicon-16x16.png
srv/http/pihole/admin/img/favicons/favicon-32x32.png
srv/http/pihole/admin/img/favicons/favicon.ico
srv/http/pihole/admin/img/favicons/manifest.json
srv/http/pihole/admin/img/favicons/mstile-150x150.png
srv/http/pihole/admin/img/favicons/safari-pinned-tab.svg
srv/http/pihole/admin/img/logo.svg
srv/http/pihole/admin/index.php
srv/http/pihole/admin/messages.php
srv/http/pihole/admin/network.php
srv/http/pihole/admin/package-lock.json
srv/http/pihole/admin/package.json
srv/http/pihole/admin/phpstan.neon.dist
srv/http/pihole/admin/queries.php
srv/http/pihole/admin/queryads.php
srv/http/pihole/admin/scripts/
srv/http/pihole/admin/scripts/pi-hole/
srv/http/pihole/admin/scripts/pi-hole/js/
srv/http/pihole/admin/scripts/pi-hole/js/auditlog.js
srv/http/pihole/admin/scripts/pi-hole/js/customcname.js
srv/http/pihole/admin/scripts/pi-hole/js/customdns.js
srv/http/pihole/admin/scripts/pi-hole/js/db_graph.js
srv/http/pihole/admin/scripts/pi-hole/js/db_lists.js
srv/http/pihole/admin/scripts/pi-hole/js/db_queries.js
srv/http/pihole/admin/scripts/pi-hole/js/debug.js
srv/http/pihole/admin/scripts/pi-hole/js/footer.js
srv/http/pihole/admin/scripts/pi-hole/js/gravity.js
srv/http/pihole/admin/scripts/pi-hole/js/groups-adlists.js
srv/http/pihole/admin/scripts/pi-hole/js/groups-clients.js
srv/http/pihole/admin/scripts/pi-hole/js/groups-domains.js
srv/http/pihole/admin/scripts/pi-hole/js/groups.js
srv/http/pihole/admin/scripts/pi-hole/js/index.js
srv/http/pihole/admin/scripts/pi-hole/js/ip-address-sorting.js
srv/http/pihole/admin/scripts/pi-hole/js/messages.js
srv/http/pihole/admin/scripts/pi-hole/js/network.js
srv/http/pihole/admin/scripts/pi-hole/js/queries.js
srv/http/pihole/admin/scripts/pi-hole/js/queryads.js
srv/http/pihole/admin/scripts/pi-hole/js/settings.js
srv/http/pihole/admin/scripts/pi-hole/js/taillog-FTL.js
srv/http/pihole/admin/scripts/pi-hole/js/taillog.js
srv/http/pihole/admin/scripts/pi-hole/js/utils.js
srv/http/pihole/admin/scripts/pi-hole/php/
srv/http/pihole/admin/scripts/pi-hole/php/FTL.php
srv/http/pihole/admin/scripts/pi-hole/php/api_token.php
srv/http/pihole/admin/scripts/pi-hole/php/auth.php
srv/http/pihole/admin/scripts/pi-hole/php/customcname.php
srv/http/pihole/admin/scripts/pi-hole/php/customdns.php
srv/http/pihole/admin/scripts/pi-hole/php/database.php
srv/http/pihole/admin/scripts/pi-hole/php/debug.php
srv/http/pihole/admin/scripts/pi-hole/php/footer.php
srv/http/pihole/admin/scripts/pi-hole/php/func.php
srv/http/pihole/admin/scripts/pi-hole/php/gravity.php
srv/http/pihole/admin/scripts/pi-hole/php/gravity.sh.php
srv/http/pihole/admin/scripts/pi-hole/php/groups.php
srv/http/pihole/admin/scripts/pi-hole/php/header.php
srv/http/pihole/admin/scripts/pi-hole/php/loginpage.php
srv/http/pihole/admin/scripts/pi-hole/php/message.php
srv/http/pihole/admin/scripts/pi-hole/php/network.php
srv/http/pihole/admin/scripts/pi-hole/php/password.php
srv/http/pihole/admin/scripts/pi-hole/php/queryads.php
srv/http/pihole/admin/scripts/pi-hole/php/savesettings.php
srv/http/pihole/admin/scripts/pi-hole/php/tailLog.php
srv/http/pihole/admin/scripts/pi-hole/php/teleporter.php
srv/http/pihole/admin/scripts/pi-hole/php/theme.php
srv/http/pihole/admin/scripts/pi-hole/php/update_checker.php
srv/http/pihole/admin/scripts/vendor/
srv/http/pihole/admin/scripts/vendor/Chart.min.js
srv/http/pihole/admin/scripts/vendor/adminlte.min.js
srv/http/pihole/admin/scripts/vendor/bootstrap-notify.min.js
srv/http/pihole/admin/scripts/vendor/bootstrap-select.min.js
srv/http/pihole/admin/scripts/vendor/bootstrap-select.min.js.map
srv/http/pihole/admin/scripts/vendor/bootstrap-toggle.min.js
srv/http/pihole/admin/scripts/vendor/bootstrap-toggle.min.js.map
srv/http/pihole/admin/scripts/vendor/datatables.min.js
srv/http/pihole/admin/scripts/vendor/daterangepicker.min.js
srv/http/pihole/admin/scripts/vendor/daterangepicker.min.js.map
srv/http/pihole/admin/scripts/vendor/jquery.confirm.min.js
srv/http/pihole/admin/scripts/vendor/jquery.min.js
srv/http/pihole/admin/scripts/vendor/moment.min.js
srv/http/pihole/admin/scripts/vendor/qrcode.php
srv/http/pihole/admin/scripts/vendor/select2.min.js
srv/http/pihole/admin/settings.php
srv/http/pihole/admin/style/
srv/http/pihole/admin/style/pi-hole.css
srv/http/pihole/admin/style/themes/
srv/http/pihole/admin/style/themes/default-auto.css
srv/http/pihole/admin/style/themes/default-dark.css
srv/http/pihole/admin/style/themes/default-darker.css
srv/http/pihole/admin/style/themes/default-light.css
srv/http/pihole/admin/style/themes/lcars.css
srv/http/pihole/admin/style/vendor/
srv/http/pihole/admin/style/vendor/AdminLTE.min.css
srv/http/pihole/admin/style/vendor/SourceSansPro/
srv/http/pihole/admin/style/vendor/SourceSansPro/SourceSansPro.css
srv/http/pihole/admin/style/vendor/SourceSansPro/source-sans-pro-v13-latin-300.woff
srv/http/pihole/admin/style/vendor/SourceSansPro/source-sans-pro-v13-latin-300.woff2
srv/http/pihole/admin/style/vendor/SourceSansPro/source-sans-pro-v13-latin-300italic.woff
srv/http/pihole/admin/style/vendor/SourceSansPro/source-sans-pro-v13-latin-300italic.woff2
srv/http/pihole/admin/style/vendor/SourceSansPro/source-sans-pro-v13-latin-600.woff
srv/http/pihole/admin/style/vendor/SourceSansPro/source-sans-pro-v13-latin-600.woff2
srv/http/pihole/admin/style/vendor/SourceSansPro/source-sans-pro-v13-latin-600italic.woff
srv/http/pihole/admin/style/vendor/SourceSansPro/source-sans-pro-v13-latin-600italic.woff2
srv/http/pihole/admin/style/vendor/SourceSansPro/source-sans-pro-v13-latin-700.woff
srv/http/pihole/admin/style/vendor/SourceSansPro/source-sans-pro-v13-latin-700.woff2
srv/http/pihole/admin/style/vendor/SourceSansPro/source-sans-pro-v13-latin-italic.woff
srv/http/pihole/admin/style/vendor/SourceSansPro/source-sans-pro-v13-latin-italic.woff2
srv/http/pihole/admin/style/vendor/SourceSansPro/source-sans-pro-v13-latin-regular.woff
srv/http/pihole/admin/style/vendor/SourceSansPro/source-sans-pro-v13-latin-regular.woff2
srv/http/pihole/admin/style/vendor/animate.min.css
srv/http/pihole/admin/style/vendor/bootstrap-select.min.css
srv/http/pihole/admin/style/vendor/bootstrap-toggle.min.css
srv/http/pihole/admin/style/vendor/bootstrap/
srv/http/pihole/admin/style/vendor/bootstrap/css/
srv/http/pihole/admin/style/vendor/bootstrap/css/bootstrap.min.css
srv/http/pihole/admin/style/vendor/bootstrap/css/bootstrap.min.css.map
srv/http/pihole/admin/style/vendor/bootstrap/fonts/
srv/http/pihole/admin/style/vendor/bootstrap/fonts/glyphicons-halflings-regular.eot
srv/http/pihole/admin/style/vendor/bootstrap/fonts/glyphicons-halflings-regular.svg
srv/http/pihole/admin/style/vendor/bootstrap/fonts/glyphicons-halflings-regular.ttf
srv/http/pihole/admin/style/vendor/bootstrap/fonts/glyphicons-halflings-regular.woff
srv/http/pihole/admin/style/vendor/bootstrap/fonts/glyphicons-halflings-regular.woff2
srv/http/pihole/admin/style/vendor/bootstrap/js/
srv/http/pihole/admin/style/vendor/bootstrap/js/bootstrap.min.js
srv/http/pihole/admin/style/vendor/datatables.min.css
srv/http/pihole/admin/style/vendor/daterangepicker.min.css
srv/http/pihole/admin/style/vendor/font-awesome/
srv/http/pihole/admin/style/vendor/font-awesome/js/
srv/http/pihole/admin/style/vendor/font-awesome/js/all.min.js
srv/http/pihole/admin/style/vendor/fonts/
srv/http/pihole/admin/style/vendor/fonts/antonio/
srv/http/pihole/admin/style/vendor/fonts/antonio/antonio-v1-latin-ext_latin-100.woff
srv/http/pihole/admin/style/vendor/fonts/antonio/antonio-v1-latin-ext_latin-100.woff2
srv/http/pihole/admin/style/vendor/fonts/antonio/antonio-v1-latin-ext_latin-700.woff
srv/http/pihole/admin/style/vendor/fonts/antonio/antonio-v1-latin-ext_latin-700.woff2
srv/http/pihole/admin/style/vendor/fonts/antonio/antonio-v1-latin-ext_latin-regular.woff
srv/http/pihole/admin/style/vendor/fonts/antonio/antonio-v1-latin-ext_latin-regular.woff2
srv/http/pihole/admin/style/vendor/fonts/antonio/antonio.css
srv/http/pihole/admin/style/vendor/fonts/ubuntu-mono/
srv/http/pihole/admin/style/vendor/fonts/ubuntu-mono/ubuntu-mono-v10-latin-ext_latin-700.woff
srv/http/pihole/admin/style/vendor/fonts/ubuntu-mono/ubuntu-mono-v10-latin-ext_latin-700.woff2
srv/http/pihole/admin/style/vendor/fonts/ubuntu-mono/ubuntu-mono-v10-latin-ext_latin-700italic.woff
srv/http/pihole/admin/style/vendor/fonts/ubuntu-mono/ubuntu-mono-v10-latin-ext_latin-700italic.woff2
srv/http/pihole/admin/style/vendor/fonts/ubuntu-mono/ubuntu-mono-v10-latin-ext_latin-italic.woff
srv/http/pihole/admin/style/vendor/fonts/ubuntu-mono/ubuntu-mono-v10-latin-ext_latin-italic.woff2
srv/http/pihole/admin/style/vendor/fonts/ubuntu-mono/ubuntu-mono-v10-latin-ext_latin-regular.woff
srv/http/pihole/admin/style/vendor/fonts/ubuntu-mono/ubuntu-mono-v10-latin-ext_latin-regular.woff2
srv/http/pihole/admin/style/vendor/fonts/ubuntu-mono/ubuntu-mono.css
srv/http/pihole/admin/style/vendor/icheck-bootstrap.min.css
srv/http/pihole/admin/style/vendor/icheck-material.min.css
srv/http/pihole/admin/style/vendor/js-warn.css
srv/http/pihole/admin/style/vendor/select2.min.css
srv/http/pihole/admin/taillog-FTL.php
srv/http/pihole/admin/taillog.php
srv/http/pihole/pihole/
srv/http/pihole/pihole/blockingpage.css
srv/http/pihole/pihole/index.php
usr/
usr/bin/
usr/bin/pihole
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/multi-user.target.wants/
usr/lib/systemd/system/multi-user.target.wants/pi-hole-gravity.timer
usr/lib/systemd/system/multi-user.target.wants/pi-hole-logtruncate.timer
usr/lib/systemd/system/pi-hole-gravity.service
usr/lib/systemd/system/pi-hole-gravity.timer
usr/lib/systemd/system/pi-hole-logtruncate.service
usr/lib/systemd/system/pi-hole-logtruncate.timer
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/pi-hole.conf
usr/share/
usr/share/licenses/
usr/share/licenses/pihole/
usr/share/licenses/pihole/AdminLTE
usr/share/licenses/pihole/Pi-hole
usr/share/pihole/
usr/share/pihole/configs/
usr/share/pihole/configs/lighttpd.example.conf
usr/share/pihole/configs/nginx.example.conf
