# pi-hole-ftl-5.14-1/files
%FILES%
etc/
etc/pihole/
etc/pihole/pihole-FTL.conf
etc/pihole/pihole-FTL.db
usr/
usr/bin/
usr/bin/dnsmasq
usr/bin/pihole-FTL
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/multi-user.target.wants/
usr/lib/systemd/system/multi-user.target.wants/pihole-FTL.service
usr/lib/systemd/system/pihole-FTL.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/pi-hole-ftl.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/pi-hole-ftl.conf
usr/share/
usr/share/licenses/
usr/share/licenses/pihole/
usr/share/licenses/pihole/Pi-hole-FTL
