# mozillavpn-2.7.0-3/files
%FILES%
etc/
etc/chromium/
etc/chromium/native-messaging-hosts/
etc/chromium/native-messaging-hosts/mozillavpn.json
etc/opt/
etc/opt/chrome/
etc/opt/chrome/native-messaging-hosts/
etc/opt/chrome/native-messaging-hosts/mozillavpn.json
etc/xdg/
etc/xdg/autostart/
etc/xdg/autostart/MozillaVPN-startup.desktop
usr/
usr/bin/
usr/bin/mozillavpn
usr/lib/
usr/lib/mozilla/
usr/lib/mozilla/native-messaging-hosts/
usr/lib/mozilla/native-messaging-hosts/mozillavpn.json
usr/lib/mozillavpn/
usr/lib/mozillavpn/mozillavpnnp
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/mozillavpn.service
usr/share/
usr/share/applications/
usr/share/applications/MozillaVPN.desktop
usr/share/dbus-1/
usr/share/dbus-1/system-services/
usr/share/dbus-1/system-services/org.mozilla.vpn.dbus.service
usr/share/dbus-1/system.d/
usr/share/dbus-1/system.d/org.mozilla.vpn.conf
usr/share/icons/
usr/share/icons/hicolor/
usr/share/icons/hicolor/128x128/
usr/share/icons/hicolor/128x128/apps/
usr/share/icons/hicolor/128x128/apps/mozillavpn.png
usr/share/icons/hicolor/16x16/
usr/share/icons/hicolor/16x16/apps/
usr/share/icons/hicolor/16x16/apps/mozillavpn.png
usr/share/icons/hicolor/32x32/
usr/share/icons/hicolor/32x32/apps/
usr/share/icons/hicolor/32x32/apps/mozillavpn.png
usr/share/icons/hicolor/48x48/
usr/share/icons/hicolor/48x48/apps/
usr/share/icons/hicolor/48x48/apps/mozillavpn.png
usr/share/icons/hicolor/64x64/
usr/share/icons/hicolor/64x64/apps/
usr/share/icons/hicolor/64x64/apps/mozillavpn.png
usr/share/polkit-1/
usr/share/polkit-1/actions/
usr/share/polkit-1/actions/org.mozilla.vpn.policy
