# photoview-2.3.12-1/files
%FILES%
etc/
etc/photoview.env
usr/
usr/bin/
usr/bin/photoview
usr/lib/
usr/lib/photoview/
usr/lib/photoview/data/
usr/lib/photoview/data/models/
usr/lib/photoview/data/models/dlib_face_recognition_resnet_model_v1.dat
usr/lib/photoview/data/models/mmod_human_face_detector.dat
usr/lib/photoview/data/models/shape_predictor_5_face_landmarks.dat
usr/lib/photoview/photoview
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/photoview.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/photoview.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/photoview.conf
usr/share/
usr/share/webapps/
usr/share/webapps/photoview-ui/
usr/share/webapps/photoview-ui/apple-touch-icon.png
usr/share/webapps/photoview-ui/asset-manifest.json
usr/share/webapps/photoview-ui/favicon.ico
usr/share/webapps/photoview-ui/index.html
usr/share/webapps/photoview-ui/logo192.png
usr/share/webapps/photoview-ui/logo512.png
usr/share/webapps/photoview-ui/manifest.json
usr/share/webapps/photoview-ui/photoview-logo.svg
usr/share/webapps/photoview-ui/service-worker.js
usr/share/webapps/photoview-ui/service-worker.js.map
usr/share/webapps/photoview-ui/static/
usr/share/webapps/photoview-ui/static/css/
usr/share/webapps/photoview-ui/static/css/1.76eb4979.chunk.css
usr/share/webapps/photoview-ui/static/css/1.76eb4979.chunk.css.map
usr/share/webapps/photoview-ui/static/css/main.53033c04.chunk.css
usr/share/webapps/photoview-ui/static/css/main.53033c04.chunk.css.map
usr/share/webapps/photoview-ui/static/js/
usr/share/webapps/photoview-ui/static/js/0.fea3a025.chunk.js
usr/share/webapps/photoview-ui/static/js/0.fea3a025.chunk.js.map
usr/share/webapps/photoview-ui/static/js/1.a2e56b8f.chunk.js
usr/share/webapps/photoview-ui/static/js/1.a2e56b8f.chunk.js.map
usr/share/webapps/photoview-ui/static/js/10.1b42e9da.chunk.js
usr/share/webapps/photoview-ui/static/js/10.1b42e9da.chunk.js.map
usr/share/webapps/photoview-ui/static/js/11.7f994222.chunk.js
usr/share/webapps/photoview-ui/static/js/11.7f994222.chunk.js.map
usr/share/webapps/photoview-ui/static/js/12.1089a71d.chunk.js
usr/share/webapps/photoview-ui/static/js/12.1089a71d.chunk.js.map
usr/share/webapps/photoview-ui/static/js/13.ea135e3a.chunk.js
usr/share/webapps/photoview-ui/static/js/13.ea135e3a.chunk.js.map
usr/share/webapps/photoview-ui/static/js/14.d9e2a806.chunk.js
usr/share/webapps/photoview-ui/static/js/14.d9e2a806.chunk.js.map
usr/share/webapps/photoview-ui/static/js/15.d124d326.chunk.js
usr/share/webapps/photoview-ui/static/js/15.d124d326.chunk.js.map
usr/share/webapps/photoview-ui/static/js/16.a17d33e6.chunk.js
usr/share/webapps/photoview-ui/static/js/16.a17d33e6.chunk.js.map
usr/share/webapps/photoview-ui/static/js/17.3259e92b.chunk.js
usr/share/webapps/photoview-ui/static/js/17.3259e92b.chunk.js.map
usr/share/webapps/photoview-ui/static/js/18.88a2297b.chunk.js
usr/share/webapps/photoview-ui/static/js/18.88a2297b.chunk.js.map
usr/share/webapps/photoview-ui/static/js/19.3171267c.chunk.js
usr/share/webapps/photoview-ui/static/js/19.3171267c.chunk.js.map
usr/share/webapps/photoview-ui/static/js/2.c87fce1c.chunk.js
usr/share/webapps/photoview-ui/static/js/2.c87fce1c.chunk.js.map
usr/share/webapps/photoview-ui/static/js/20.00064f58.chunk.js
usr/share/webapps/photoview-ui/static/js/20.00064f58.chunk.js.map
usr/share/webapps/photoview-ui/static/js/21.360b0aeb.chunk.js
usr/share/webapps/photoview-ui/static/js/21.360b0aeb.chunk.js.map
usr/share/webapps/photoview-ui/static/js/22.75ebe866.chunk.js
usr/share/webapps/photoview-ui/static/js/22.75ebe866.chunk.js.map
usr/share/webapps/photoview-ui/static/js/23.a8e29a4f.chunk.js
usr/share/webapps/photoview-ui/static/js/23.a8e29a4f.chunk.js.map
usr/share/webapps/photoview-ui/static/js/24.f6f0514d.chunk.js
usr/share/webapps/photoview-ui/static/js/24.f6f0514d.chunk.js.map
usr/share/webapps/photoview-ui/static/js/25.a85db325.chunk.js
usr/share/webapps/photoview-ui/static/js/25.a85db325.chunk.js.map
usr/share/webapps/photoview-ui/static/js/26.3f306e50.chunk.js
usr/share/webapps/photoview-ui/static/js/26.3f306e50.chunk.js.map
usr/share/webapps/photoview-ui/static/js/27.ddb8cc18.chunk.js
usr/share/webapps/photoview-ui/static/js/27.ddb8cc18.chunk.js.map
usr/share/webapps/photoview-ui/static/js/28.f5a9f5ca.chunk.js
usr/share/webapps/photoview-ui/static/js/28.f5a9f5ca.chunk.js.map
usr/share/webapps/photoview-ui/static/js/29.668ed4bc.chunk.js
usr/share/webapps/photoview-ui/static/js/29.668ed4bc.chunk.js.map
usr/share/webapps/photoview-ui/static/js/3.2c1660b4.chunk.js
usr/share/webapps/photoview-ui/static/js/3.2c1660b4.chunk.js.map
usr/share/webapps/photoview-ui/static/js/4.994c78d5.chunk.js
usr/share/webapps/photoview-ui/static/js/4.994c78d5.chunk.js.map
usr/share/webapps/photoview-ui/static/js/5.4539af1c.chunk.js
usr/share/webapps/photoview-ui/static/js/5.4539af1c.chunk.js.map
usr/share/webapps/photoview-ui/static/js/8.f09a7e8c.chunk.js
usr/share/webapps/photoview-ui/static/js/8.f09a7e8c.chunk.js.LICENSE.txt
usr/share/webapps/photoview-ui/static/js/8.f09a7e8c.chunk.js.map
usr/share/webapps/photoview-ui/static/js/9.4fd50470.chunk.js
usr/share/webapps/photoview-ui/static/js/9.4fd50470.chunk.js.map
usr/share/webapps/photoview-ui/static/js/main.75ee36da.chunk.js
usr/share/webapps/photoview-ui/static/js/main.75ee36da.chunk.js.map
usr/share/webapps/photoview-ui/static/js/runtime-main.d93b8ecb.js
usr/share/webapps/photoview-ui/static/js/runtime-main.d93b8ecb.js.map
usr/share/webapps/photoview-ui/static/media/
usr/share/webapps/photoview-ui/static/media/closeSidebarIcon.9f91c3bf.svg
usr/share/webapps/photoview-ui/static/media/date.770f6900.svg
usr/share/webapps/photoview-ui/static/media/direction-arrow.fabb7a2d.svg
usr/share/webapps/photoview-ui/static/media/dismissIcon.c74dc15c.svg
usr/share/webapps/photoview-ui/static/media/gear.f290e670.svg
usr/share/webapps/photoview-ui/static/media/image-popup.c4215c79.svg
usr/share/webapps/photoview-ui/static/media/peopleDotsIcon.1a677723.svg
usr/share/webapps/photoview-ui/static/media/pinSidebarIconFilled.f240d5a2.svg
usr/share/webapps/photoview-ui/static/media/pinSidebarIconOutline.d5045b27.svg
usr/share/webapps/photoview-ui/static/media/shareAddIcon.44b64157.svg
usr/share/webapps/photoview-ui/static/media/shareCopyIcon.0b9c5972.svg
usr/share/webapps/photoview-ui/static/media/shareDeleteIcon.51ac075d.svg
usr/share/webapps/photoview-ui/static/media/shareLinkIcon.52c3973d.svg
usr/share/webapps/photoview-ui/static/media/shareMoreIcon.838cbd23.svg
usr/share/webapps/photoview-ui/static/media/sorting.37d44b61.svg
usr/share/webapps/photoview-ui/static/media/textboxActionArrow.e14ab55e.svg
usr/share/webapps/photoview-ui/static/media/textboxLoadingSpinner.682e5343.svg
usr/share/webapps/photoview-ui/static/media/videoThumbnailIcon.e4d4b5e7.svg
var/
var/cache/
var/cache/photoview/
var/cache/photoview/media_cache/
var/lib/
var/lib/photoview/
